# scripts/99-flygate-ipad.rules
# FlyGate ACI - iPad USB Detection Rules
# Install to: /etc/udev/rules.d/99-flygate-ipad.rules
#
# Apple iPad USB Vendor ID: 05ac
# iPad Product IDs vary by model:
#   12a8 - iPad (various)
#   12ab - iPad mini
#   12a4/12a5 - iPad Air
#   129e+ - iPad Pro
#
# This rule triggers on any Apple iOS device connection
# Note: For production, the FlyGate agent identifies itself via the handshake,
# not via USB attributes. The udev rule just triggers the connection flow.

# iPad connected - trigger attach (any Apple iOS device)
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="05ac", ENV{DEVTYPE}=="usb_device", TAG+="systemd", RUN+="/opt/flygate-aci/scripts/usb-monitor.sh FlyGateAgent-iPad-0001 attach"

# iPad disconnected - use unbind event which is more reliable
ACTION=="unbind", SUBSYSTEM=="usb", ENV{PRODUCT}=="5ac/*", RUN+="/opt/flygate-aci/scripts/usb-monitor.sh FlyGateAgent-iPad-0001 detach"

# Alternative remove trigger using ID_VENDOR
ACTION=="remove", SUBSYSTEM=="usb", ENV{ID_VENDOR}=="Apple*", RUN+="/opt/flygate-aci/scripts/usb-monitor.sh FlyGateAgent-iPad-0001 detach"
